# dev-python/pantalaimon/files/pantalaimon-openrc
#!/sbin/openrc-run

command="/usr/bin/pantalaimon"
pidfile="/run/${RC_SVCNAME}.pid"
command_background=true
